# reports/reports.cron
15 * * * * DIR=/home/ubuntu/tribtuary /home/ubuntu/tributary/reports/reports.sh
0 9 * * * /home/ubuntu/tributary/reports/backup.sh

